# systemd/user/ags.service
[Unit]
Description=Aylur's GTK Shell

[Service]
Type=simple
ExecStart=ags

