# systemd/user/clean-xdg-cache-home.timer
[Unit]
Description=Weekly cleanup of XDG_CACHE_HOME

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target

