# systemd/user/dunst.service
[Unit]
Description=Dunst notification daemon

[Service]
Type=simple
ExecStart=dunst

