# systemd/user/foot-popup-server.service
[Unit]
Description=Foot server for foot popups

[Service]
Type=simple
ExecStart=foot --server=/tmp/foot-popup-server.sock
Restart=always

