# systemd/user/hypridle.service
[Unit]
Description=Hyprland's idle daemon

[Service]
Type=simple
ExecStart=hypridle

