# systemd/user/hyprpaper.service
[Unit]
Description=Hyprland's wallpaper daemon

[Service]
Type=simple
ExecStart=hyprpaper

