# systemd/user/xdg-desktop-portal.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/libexec/xdg-desktop-portal --verbose
